# etc/systemd/system/boot-update.service
[Unit]
Description=Update BootC Container on Boot
Wants=network-online.target
After=sshd.service
ConditionPathExists=/run/ostree-booted

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc update --apply --quiet
RemainAfterExit=true

[Install]
WantedBy=multi-user.target