# etc/systemd/system/boot-update.service
[Unit]
Description=Update on boot.
After=network-online.target
Before=multi-user.target reboot.target default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bootc update --apply --quiet

[Install]
WantedBy=multi-user.target default.target
